# etc/init/deer_site.conf
description "deer mailbox recorder"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
script
	exec /home/pi/deer_camera/www/deer_site.py
end script
