# etc/init/deer.conf
description "deer mailbox recorder"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
script
	exec /home/pi/deer_camera/src/deer.py /home/pi/deer_camera/cfg/deer.cfg
end script
